# packages/systemd/user/backup_notes.service
[Unit]
Description=Backup notes

[Service]
ExecStart=/home/damien/.zsh/scripts/backup-notes
