# packages/systemd/user/notify_updates.timer
[Unit]
Description=Notify about packages updates

[Timer]
OnBootSec=1h
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
