# packages/systemd/user/notify_package_updates.service
[Unit]
Description=Notify about packages updates

[Service]
Type=oneshot
ExecStart=/home/damien/.config/sway/scripts/check_notify_package_updates
