# packages/systemd/user/notify_updates.service
[Unit]
Description=Notify about packages updates

[Service]
Type=oneshot
ExecStart=/home/damien/.config/sway/scripts/notify_updates
