# packages/systemd/user/notify_low_battery.service
[Unit]
Description=Notify about low battery

[Service]
Type=oneshot
ExecStart=/home/damien/.config/sway/scripts/notify_low_battery
