# packages/systemd/user/qutebrowser_update_adblock.service
[Unit]
Description=Update qutebrowser adblock list

[Service]
Type=oneshot
ExecStart=/home/damien/.config/qutebrowser/scripts/adblock_update
