# packages/systemd/user/qutebrowser_update_adblock.timer
[Unit]
Description=Update qutebrowser adblock list

[Timer]
OnBootSec=4h

[Install]
WantedBy=timers.target
