# packages/systemd/user/qutebrowser_update_adblock.timer
[Unit]
Description=Update qutebrowser adblock list

[Timer]
OnCalendar=weekly
Persistent=True

[Install]
WantedBy=timers.target
