# systemd/user/clear_do_not_disturb.service
[Unit]
Description=Clear "Do not disturb" mode
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=/home/damien/.config/sway/scripts/clear_do_not_disturb
