# systemd/user/notify_fortune.timer
[Unit]
Description=Notify about the fortune of the day

[Timer]
OnBootSec=5m

[Install]
WantedBy=timers.target
