# configs/systemd/user/backup_bookmarks.service
[Unit]
Description=Backup bookmarks

[Service]
ExecStart={dir/bin}/backup-bookmarks
