# configs/systemd/user/backup_notes.service
[Unit]
Description=Backup notes

[Service]
ExecStart=/home/damien/.local/bin/backup-notes
