# configs/systemd/user/backup_bookmarks.service
[Unit]
Description=Backup bookmarks

[Service]
ExecStart={folder/bin}/backup-bookmarks
