# configs/systemd/user/backup_notes.service
[Unit]
Description=Backup notes

[Service]
ExecStart={folder/bin}/backup-notes
