# configs/systemd/user/backup_notes.service
[Unit]
Description=Backup notes

[Service]
ExecStart={dir/bin}/backup-notes
