# configs/systemd/user/backup_bookmarks.path
[Unit]
Description=Backup bookmarks

[Path]
PathChanged=/home/damien/.config/qutebrowser/bookmarks/urls
PathChanged=/home/damien/.config/qutebrowser/quickmarks

[Install]
WantedBy=paths.target
