# configs/systemd/user/backup_notes.path
[Unit]
Description=Backup notes

[Path]
PathChanged=/home/damien/Notes

[Install]
WantedBy=paths.target
