# configs/systemd/user/backup_bookmarks.service
[Unit]
Description=Backup bookmarks

[Service]
ExecStart=/home/damien/.zsh/scripts/backup-bookmarks
