# configs/systemd/user/notify_package_updates.service
[Unit]
Description=Notify about packages updates

[Service]
Type=oneshot
ExecStart={folder/bin}/check-notify-package-updates
