# configs/systemd/user/notify_low_battery.service
[Unit]
Description=Notify about low battery

[Service]
Type=oneshot
ExecStart={dir/bin}/check-notify-low-battery
