# configs/systemd/user/notify_package_updates.service
[Unit]
Description=Notify about packages updates

[Service]
Type=oneshot
ExecStart=/home/damien/.zsh/scripts/check-notify-package-updates
