# configs/systemd/user/qutebrowser_update_adblock.service
[Unit]
Description=Update qutebrowser adblock list

[Service]
Type=oneshot
ExecStart=/usr/bin/qutebrowser :adblock-update
