# files/configs/systemd/user/backup_archives.path
[Unit]
Description=Backup archives

[Path]
PathChanged=%h/Archives

[Install]
WantedBy=paths.target
