# files/configs/systemd/user/backup_todo.timer
[Unit]
Description=Backup TODO files

[Timer]
OnStartupSec=2H
OnUnitActiveSec=2H

[Install]
WantedBy=timers.target
