# files/configs/systemd/user/backup_contacts.service
[Unit]
Description=Backup contacts

[Service]
Environment=SSH_AUTH_SOCK=%t/gnupg/S.gpg-agent.ssh
ExecStart={dir/bin}/backup-contacts
