# files/configs/systemd/user/backup_time_tracker_logs.path
[Unit]
Description=Backup time tracker logs

[Path]
PathChanged={dir/data}/time-tracker/logs

[Install]
WantedBy=paths.target
