# files/configs/systemd/user/backup_notes.path
[Unit]
Description=Backup notes

[Path]
PathChanged=%h/Notes

[Install]
WantedBy=paths.target
