# files/configs/systemd/user/backup_task_logs.path
[Unit]
Description=Backup task logs

[Path]
PathChanged={dir/data}/task

[Install]
WantedBy=paths.target
