# files/configs/systemd/user/backup_todo.timer
[Unit]
Description=Backup TODO files

[Timer]
OnStartupSec=2h
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
