# files/configs/systemd/user/backup_todo.path
[Unit]
Description=Backup TODO files

[Path]
PathChanged=%h/Todo

[Install]
WantedBy=paths.target
