# files/configs/systemd/user/backup_task_logs.service
[Unit]
Description=Backup task logs

[Service]
Environment=SSH_AUTH_SOCK=%t/gnupg/S.gpg-agent.ssh
ExecStart={dir/bin}/backup task_logs
