# files/configs/systemd/user/documents_build_cache.timer
[Unit]
Description=Build documents list cache

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
