# files/configs/systemd/user/documents_build_cache.service
[Unit]
Description=Build documents list cache

[Service]
Type=oneshot
ExecStart={dir/bin}/open-document -c
