# files/configs/systemd/user/documents_build_cache.timer
[Unit]
Description=Build documents list cache

[Timer]
OnStartupSec=1h

[Install]
WantedBy=timers.target
