# files/configs/systemd/user/documentation_build_cache.service
[Unit]
Description=Build documentation topics list cache

[Service]
Type=oneshot
ExecStart={dir/bin}/documentation-all -c
