# files/configs/systemd/user/documentation_build_cache.timer
[Unit]
Description=Build documentation topics list cache

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
