# files/configs/systemd/user/documentation_build_cache.timer
[Unit]
Description=Build documentation topics list cache

[Timer]
OnStartupSec=1h

[Install]
WantedBy=timers.target
