# files/configs/systemd/user/fetch_news.service
[Unit]
Description=Fetch news

[Service]
Type=oneshot
ExecStart=/usr/bin/sfeed_update
