# files/configs/systemd/user/fetch_news.timer
[Unit]
Description=Fetch news

[Timer]
OnStartupSec=10m
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
