# files/configs/systemd/user/notify_break.timer
[Unit]
Description=Notify break time

[Timer]
OnBootSec=20min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
