# files/configs/systemd/user/notify_low_battery.timer
[Unit]
Description=Notify about low battery

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
