# files/configs/systemd/user/notify_break.timer
[Unit]
Description=Notify break time

[Timer]
OnStartupSec=20min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
