# files/configs/systemd/user/notify_break.service
[Unit]
Description=Notify break time

[Service]
Type=oneshot
ExecStart={dir/bin}/notify-break
