# files/configs/systemd/user/qutebrowser_update_adblock.timer
[Unit]
Description=Update qutebrowser adblock list

[Timer]
OnStartupSec=4h

[Install]
WantedBy=timers.target
