# files/configs/systemd/user/sync_mails.service
[Unit]
Description=Synchronize mails

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va
