# files/configs/systemd/user/sync_mails.service
[Unit]
Description=Synchronize mails

[Service]
ExecStart=/usr/bin/mbsync -a
