# files/configs/systemd/user/sync_mails.timer
[Unit]
Description=Synchronize mails

[Timer]
OnStartupSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
