# files/configs/systemd/user/vdirsyncer.service
[Unit]
Description=Synchronize calendars and contacts

[Service]
ExecStart=/usr/bin/vdirsyncer sync
RuntimeMaxSec=3m
