# files/configs/systemd/user/vdirsyncer.timer
[Unit]
Description=Synchronize vdirs

[Timer]
OnStartupSec=1h

[Install]
WantedBy=timers.target
